# scripts/softap/cbridge-setup-ap.service
[Unit]
Description=C-Bridge Setup Access Point Controller
Documentation=file:///home/user1/C-Bridge-Production/scripts/softap/README.md
After=network.target NetworkManager.service
Wants=network.target

[Service]
Type=simple
# Path will be updated by setup_script during installation to actual location
# Placeholder - setup_script will replace this with actual path
ExecStart=/home/user1/C-Bridge-Production/scripts/softap/auto-ap.sh run
# Continuously monitor WiFi connection and internet connectivity
# Automatically starts SoftAP when WiFi client mode disconnects or loses internet
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

